# rootfs/etc/systemd/system/r8q-usb-gadget.service
[Unit]
Description=r8q USB NCM gadget (ssh over USB)
After=multi-user.target
[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/sbin/r8q-usb-gadget-up.sh
[Install]
WantedBy=multi-user.target
