# rootfs/etc/systemd/system/r8q-battery.service
[Unit]
Description=r8q battery fuel gauge (MAX77705 ModelGauge @ SE0) module load
# The geni i2c bus (SE0) that carries the fuel gauge is brought up by
# r8q-touch.service; the 0-0036 i2c device is auto-created from the DT
# once that binds, so just load the driver afterwards. Cold-load once,
# module lives outside /lib/modules (udev coldplug stays disabled).
After=r8q-touch.service
Requires=r8q-touch.service
ConditionPathExists=!/sys/class/power_supply/max170xx_battery

[Service]
Type=oneshot
ExecStart=/usr/bin/insmod /root/ts/max17042_battery.ko
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
