# rootfs/etc/systemd/system/gdm.service.d/r8q-after-gpu.conf
[Unit]
Wants=r8q-gpu.service
After=r8q-gpu.service
